# deploy/nidin-bos-scheduler.service
# systemd service for background scheduler (sync integrations, morning briefing)
# Place at: /etc/systemd/system/nidin-bos-scheduler.service
#
# Runs as a single process — handles periodic sync for all orgs.
# Separate from Gunicorn so multi-worker doesn't duplicate sync tasks.

[Unit]
Description=Nidin BOS Scheduler
After=network.target postgresql.service redis.service
PartOf=nidin-bos.service

[Service]
Type=notify
NotifyAccess=main
User=deploy
Group=deploy
WorkingDirectory=/opt/nidin-bos
Environment="PATH=/opt/nidin-bos/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/opt/nidin-bos/.env
ExecStart=/opt/nidin-bos/venv/bin/python run_scheduler.py
Restart=always
RestartSec=10
WatchdogSec=120
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nidin-bos
PrivateTmp=true

[Install]
WantedBy=multi-user.target
