# deploy/nidin-bos-scheduler.service
# systemd service for background scheduler (sync integrations, morning briefing)
# Place at: /etc/systemd/system/nidin-nover-ai-scheduler.service
#
# Runs as a single process — handles periodic sync for all orgs.
# Separate from Gunicorn so multi-worker doesn't duplicate sync tasks.

[Unit]
Description=Nidin Nover AI Scheduler
After=network.target postgresql.service redis.service
PartOf=nidin-nover-ai.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/nidin-nover-ai
Environment="PATH=/opt/nidin-nover-ai/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/opt/nidin-nover-ai/.env
ExecStart=/opt/nidin-nover-ai/venv/bin/python run_scheduler.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nidin-nover-ai
PrivateTmp=true

[Install]
WantedBy=multi-user.target
