# deploy/nidin-bos.service
# systemd service for Nidin BOS web server (Gunicorn + Uvicorn workers)
# Place at: /etc/systemd/system/nidin-bos.service
#
# Commands:
#   sudo systemctl enable nidin-bos
#   sudo systemctl start nidin-bos
#   sudo systemctl status nidin-bos
#   sudo journalctl -u nidin-bos -f

[Unit]
Description=Nidin BOS (Gunicorn)
After=network.target postgresql.service redis.service
Wants=nidin-bos-scheduler.service

[Service]
Type=notify
NotifyAccess=all
WatchdogSec=120
User=deploy
Group=deploy
WorkingDirectory=/opt/nidin-bos
Environment="PATH=/opt/nidin-bos/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/opt/nidin-bos/.env
ExecStart=/opt/nidin-bos/venv/bin/gunicorn -c gunicorn.conf.py app.main:app
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nidin-bos
PrivateTmp=true

[Install]
WantedBy=multi-user.target
